# udev/95-monitor-hotplug.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/edma2/bin/external-monitor"
